# ops/option-learner-guide.service
[Unit]
Description=Option Learner Guide - Next.js Application
After=network.target
Wants=network.target

[Service]
Type=simple
User=kunkka
Group=kunkka
WorkingDirectory=/home/kunkka/projects/option-learner-guide
Environment=NODE_ENV=production
Environment=PORT=3601
Environment=NEXT_TELEMETRY_DISABLED=1
ExecStart=/usr/bin/node /home/kunkka/projects/option-learner-guide/node_modules/.bin/next start -p 3601
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=option-learner-guide
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target